# .config/systemd/user/dunst.service.d/override.conf
[Unit]
After=display-manager.service

[Service]
Restart=always