# .config/systemd/user/dunst.service.d/override.conf
[Service]
Restart=always
# This is to ensure that super-early notifications don't start dunst prematurely.
ConditionEnvironment=DISPLAY